# clip_pipe/supervisor/xxai-clip.ini
[program:xxai-clip]
command=$EXE
loglevel=warn
autostart=true
autorestart=true
startretries=10
stopwaitsecs=1
startsecs=3
stdout_logfile=/var/log/supervisor/xxai-clip.log
stderr_logfile=/var/log/supervisor/xxai-clip.err.log
